# .conf/supervisor-app.conf
[program:uwsgi]
command = uwsgi --ini /etc/uwsgi/sites/mysite.ini

[program:nginx]
command = nginx